# resources/openrc/kratanet
#!/sbin/openrc-run
description="Krata Networking Daemon"
command="/usr/sbin/kratanet"
supervisor="supervise-daemon"
output_log="/var/log/kratanet.log"
error_log="/var/log/kratanet.err"

depend() {
  use kratad
}

export RUST_LOG=info
