# resources/openrc/kratad
#!/sbin/openrc-run
description="Krata Isolation Engine"
command="/usr/sbin/kratad"
supervisor="supervise-daemon"
output_log="/var/log/kratad.log"
error_log="/var/log/kratad.err"

depend() {
  use xenstored
}

export RUST_LOG=info
