# resources/systemd/kratad.service
[Unit]
Description=Krata Isolation Engine

[Service]
Restart=on-failure
Type=simple
ExecStart=/usr/sbin/kratad -l unix:///var/lib/krata/daemon.socket
Environment=RUST_LOG=info
User=root

[Install]
WantedBy=multi-user.target
